# debian/cron.d
#
# cron jobs for xivo-popc-stats
#

25 3 * * * root /usr/bin/xivo-sync
